# getweather.service
[Unit]
Description=Weather Reminder


[Service]
ExecStart=/home/vagrant/Documents/week10/check_weather

[Install]
WantedBy=multi-user.target

